# dock/kubler/images/redis/etc/service/redis/run
#!/bin/sh

REDIS_CONF=${REDIS_CONF:-/etc/redis.conf}
exec /usr/sbin/redis-server "${REDIS_CONF}" 2>&1
