# dock/gentoobb/images/nginx-proxy-conf/etc/service/docker-gen/run
#!/bin/sh
exec docker-gen -watch -only-exposed /config/template/nginx.tmpl /etc/nginx/sites-enabled/default.conf
