# dock/gentoobb/images/elasticsearch/etc/service/elasticsearch/run
#!/bin/sh

ES_HOME=${ES_HOME:="/usr/share/elasticsearch"}
ES_USER=${ES_USER:="elasticsearch"}
ES_INSTANCE=local

PIDFILE="/run/elasticsearch/elasticsearch.${ES_INSTANCE}.pid"
ES_BASE_PATH="/var/lib/elasticsearch/${ES_INSTANCE}"
ES_CONF_PATH="/etc/elasticsearch/${ES_INSTANCE}"
ES_SCRIPT_PATH="/etc/elasticsearch/${ES_INSTANCE}/scripts"
ES_LOG_PATH="/var/log/elasticsearch/${ES_INSTANCE}"
ES_DATA_PATH="${ES_BASE_PATH}/data"

check_dirs="$(dirname ${PIDFILE}) ${ES_LOG_PATH} ${ES_SCRIPT_PATH} ${ES_DATA_PATH}"
for dir in ${check_dirs}; do
    [ ! -d "${dir}" ] && mkdir -p "${dir}"
    chown ${ES_USER} "${dir}"
done;

export ES_INCLUDE="${ES_CONF_PATH}/elasticsearch.in.sh"
export JAVA_OPTS
export ES_JAVA_OPTS
export ES_HEAP_SIZE
export ES_HEAP_NEWSIZE
export ES_DIRECT_SIZE
export ES_USE_IPV4

server_command="/usr/share/elasticsearch/bin/elasticsearch"
server_args=" -p ${PIDFILE} -E default.path.home=${ES_HOME} -E default.path.conf=${ES_CONF_PATH} -E default.path.data=${ES_DATA_PATH} -E default.path.logs=${ES_LOG_PATH}"
if [ -n "${ES_MAX_FD}" ]; then
    ulimit -n ${ES_MAX_FD}
fi

exec 2>&1
su-exec elasticsearch:elasticsearch ${server_command} ${server_args}
