# bb-dock/nginx-php5.5/etc/service/php-fpm/run
#!/bin/sh
exec 2>&1
exec /usr/bin/php-fpm5.5 -y /etc/php/fpm-php5.5/php-fpm.conf -F >> /var/log/php5-fpm.log
